# _superninja_startup.conf
[program:3001_repurposebot]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python3 /workspace/repurposebot/server.py"
directory=/workspace/repurposebot
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/3001_repurposebot.err.log
stdout_logfile=/var/log/supervisor/3001_repurposebot.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
